# resources/firstboot.service
[Unit]
Description=One time boot script

[Service]
Type=oneshot
ExecStart=/firstboot.sh

[Install]
WantedBy=multi-user.target